# stream_service.service
[Unit]
Description=CosyVoice Streaming Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/CosyVoice

# ❗ 不直接跑 python
ExecStart=/home/ec2-user/CosyVoice/start.sh

Restart=always
RestartSec=10

# 防止把系统拖死
MemoryMax=12G
CPUQuota=800%
LimitNOFILE=65536

Environment=PYTHONUNBUFFERED=1

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
