# stream_service.service
[Unit]
Description=CosyVoice Streaming Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/CosyVoice

ExecStart=/home/ec2-user/CosyVoice/start.sh

# ✅ 统一使用 journald 管理日志，避免日志文件无限增长
StandardOutput=journal
StandardError=journal

# Python 行为
Environment=PYTHONUNBUFFERED=1

# ✅ 资源护栏（32GB 总内存，4 worker 模式）
MemoryMax=28G
MemoryHigh=24G
TasksMax=256
LimitNOFILE=65536

# ✅ CPU 限制，防止占满CPU
CPUQuota=400%

# ✅ 更安全的重启策略
Restart=on-failure
RestartSec=5

# ✅ 防止 stop 卡死
TimeoutStartSec=30
TimeoutStopSec=15
KillSignal=SIGTERM
KillMode=mixed

[Install]
WantedBy=multi-user.target
